# otros/autofolder.service
[Unit]
Description=Make a folder with actual date
After=network.target

[Service]
ExecStart=/bin/autofolder.sh

[Install]
WantedBy=default.target
